# scripts/athena_supervisor.conf
; Supervisor programs for Athena on the RunPod.
;
; Install by dropping this file at /etc/supervisor/conf.d/athena.conf, then
; `supervisorctl reread && supervisorctl update`.
;
; Log rotation is built in: 50MB per file, 10 files kept per stream.
; Use `supervisorctl status athena-brain athena-training` to check,
;     `supervisorctl restart athena-brain`  to bounce the brain daemon,
;     `supervisorctl tail -f athena-training stderr`  for live stderr.

[program:athena-brain]
command=python3 -u scripts/brain_daemon.py --socket /var/run/athena/brain.sock --resume --checkpoint /workspace/nimcp/checkpoints/athena/athena_immersive.bin --init-mode full
directory=/workspace/nimcp
autostart=true
autorestart=true
; Brain init is expensive (~11 min full). Allow long startup before marking RUNNING.
startsecs=30
startretries=3
; Graceful stop: SIGTERM triggers final-checkpoint save, then process exits.
stopsignal=TERM
stopwaitsecs=300
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/athena/brain.stdout
stderr_logfile=/var/log/athena/brain.stderr
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1"
user=root
priority=100

[program:athena-training]
command=python3 -u scripts/immerse_athena.py --daemon
directory=/workspace/nimcp
autostart=true
autorestart=true
; Only start training once the brain socket is ready — we use a startup gate
; inside immerse_athena.py (connect-with-retry). startsecs=30 lets it connect.
startsecs=30
startretries=5
stopsignal=TERM
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/athena/training.stdout
stderr_logfile=/var/log/athena/training.stderr
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONUNBUFFERED="1"
user=root
; Start training AFTER the brain — supervisor honors priority ordering.
priority=200
